# etc/uci-defaults/10-network
#!/bin/sh

uci set network.wan.ipaddr=172.17.0.2
uci set network.wan.gateway=172.17.0.1

uci set dhcp.@dnsmasq[0].resolvfile='/etc/resolv.conf'

uci commit 


exit 0
